# service/webserver/run
#!/bin/bash
exec 2>&1
export LANG=C
sleep 10
cd /home/lights/batshitlights/server
exec setuidgid python webserver.py
