# service/fileloop/run
#!/bin/bash
exec 2>&1
export LANG=C
sleep 0
echo $$ > /var/run/fileloop.py
cd /home/lights/batshitlights/server
exec setuidgid lights python -u filelooop.py
