# service/fileloop/run
#!/bin/bash
exec 2>&1
export LANG=C
sleep 1
cd /home/lights/batshitlights/server
exec setuidgid lights python -u filelooop.py
